# deploy/otter-letters.timer
[Unit]
Description=Run otter letters update at midnight

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target
